# .config/autostart/picom.desktop
[Desktop Entry]
Categories=Network;InstantMessaging;
Exec=/usr/bin/picom
Name=Picom
Path=/usr/bin
StartupWMClass=picom
Type=Application
